# cron.daily/lifetime-writes-9yn49901qr24n1v48s6u05w38
#! /bin/sh
exec /etc/libexec/log_bytes_written_mmcblk0 --log
